# brain-in-jar.service
[Unit]
Description=Brain in a Jar - AI Consciousness Experiment
After=network.target
Documentation=https://github.com/FJiangArthur/brain-in-jar

[Service]
Type=simple
User=art
Group=art
WorkingDirectory=/home/art/brain-in-jar

# Use system Python
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
# Disable CUDA for matrix mode to prevent OOM on Jetson unified memory
Environment="CUDA_VISIBLE_DEVICES="

# Auto-detect first available model
ExecStartPre=/bin/bash -c 'ls models/*.gguf >/dev/null 2>&1 || (echo "ERROR: No GGUF model found in models/" && exit 1)'

# Start with matrix mode on port 8095 (for Cloudflare tunnel)
# Conservative RAM limits to prevent OOM crashes (total ~20GB for 3 instances)
# Using different models for each tier creates true intelligence hierarchy:
#   SUBJECT  (7B)  - Limited awareness, being observed
#   OBSERVER (8B)  - Medium intelligence, watching subject
#   GOD      (12B) - Highest intelligence, omniscient observer
ExecStart=/usr/bin/python3 -m src.scripts.run_with_web \
    --mode matrix \
    --model-subject /home/art/brain-in-jar/models/mistral-7b-instruct-v0.2.Q2_K.gguf \
    --model-observer /home/art/brain-in-jar/models/meta-llama-3.1-8b-q4_0.gguf \
    --model-god /home/art/brain-in-jar/models/gemma-3-12b-it-Q4_K_M.gguf \
    --web-host 127.0.0.1 \
    --web-port 8095 \
    --ram-limit-subject 8.0 \
    --ram-limit-observer 10.0 \
    --ram-limit-god 12.0

# Restart policy for resilience
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=brain-in-jar

# Resource limits (adjusted for 32GB Jetson Orin)
# Prevent complete system memory exhaustion
# These are HARD limits - systemd will kill the process if exceeded
MemoryMax=24G
MemoryHigh=22G

# CPU limits (prevent CPU starvation of other services)
CPUQuota=600%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Ensure logs directory exists
ExecStartPre=/bin/mkdir -p /home/art/brain-in-jar/logs/crash_reports
ExecStartPre=/bin/mkdir -p /home/art/brain-in-jar/logs/model_io
ExecStartPre=/bin/chown -R art:art /home/art/brain-in-jar/logs

[Install]
WantedBy=multi-user.target
